# production/minfee.crontab
@reboot /usr/local/bin/namecoind >/dev/null 2>&1
